# cron/crontab
0 0,8,16 * * * root /main/run.sh >> /logs/cron.log 2>&1
